# deploy_tools/config_template/gunicorn-supervisor.conf
[program:syms]
user=USERNAME
directory=/root/sites/syms/source
command=/root/sites/syms/virtualenv/bin/gunicorn --bind unix:/tmp/syms.socket --capture-output --access-logfile ../access.log --error-logfile ../error.log --timeout=60 syms_server.wsgi:application
environment=APP_ID="APP_ID",APP_SECRET="APP_SECRET"